# templates/solr.service.erb
####################################################################
#### NOTE: THIS FILE IS PUPPET CONTROLLED - ANY CHAGES WILL BE LOST
####################################################################
[Unit]
Description=Apache SOLR
After=syslog.target network.target remote-fs.target nss-lookup.target
[Service]
PIDFile=<%= scope['::solr::solr_pid_dir'] %>/solr-<%= scope['::solr::solr_port'] %>.pid
Environment=SOLR_INCLUDE=<%= scope['::solr::solr_env'] %>
ExecStart=<%= scope['::solr::solr_bin'] -%>/solr start -noprompt
ExecRestart=<%= scope['::solr::solr_bin'] -%>/solr restart -noprompt
ExecStop=<%= scope['::solr::solr_bin'] -%>/solr stop -noprompt
ExecStatus=<%= scope['::solr::solr_bin'] -%>/solr status -noprompt
ExecReload=/bin/kill -s HUP $MAINPID
User=solr
PrivateTmp=true
[Install]
WantedBy=multi-user.target
